# Resources/LaunchAgents/com.vishruth.assistant.core.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vishruth.assistant.core</string>

    <key>MachServices</key>
    <dict>
        <key>com.vishruth.assistant.core</key>
        <true/>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>__BINARY_PATH__</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/core.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/core.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ASSISTANT_DEV</key>
        <string>1</string>
    </dict>
</dict>
</plist>
